# tools/prof/migration-scaffolding/com.mckay.prof-snapshot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- IX — prof_snapshot daily 5am, talking to straylight via tunnel. -->
<plist version="1.0">
<dict>
    <key>Label</key><string>com.mckay.prof-snapshot</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/mckay/i446-monorepo/tools/prof/prof_snapshot.py</string>
        <string>--today</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key><integer>5</integer>
        <key>Minute</key><integer>0</integer>
    </dict>
    <key>RunAtLoad</key><false/>
    <key>StandardOutPath</key><string>/Users/mckay/.config/prof/snapshot.log</string>
    <key>StandardErrorPath</key><string>/Users/mckay/.config/prof/snapshot.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key><string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>AGENCY_REMOTE_CALENDAR_PORT</key><string>7001</string>
        <key>AGENCY_REMOTE_HOST</key><string>localhost</string>
    </dict>
</dict>
</plist>
